# saas/vm-image/supervisor/programs/caddy.conf
[program:caddy]
command=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/neurolink/caddy.log
stderr_logfile=/var/log/neurolink/caddy.err.log
priority=50
