# saas/vm-image/supervisor/programs/cloudflared.conf
[program:cloudflared]
; Tunnel token is written to /secrets/cloudflared-token by first-boot.sh
; after the provisioner has called Cloudflare's API.
command=/bin/sh -c 'exec /usr/local/bin/cloudflared tunnel --no-autoupdate run --token "$(cat /secrets/cloudflared-token)"'
autostart=true
autorestart=true
; Wait for Caddy (50) to be up before the tunnel fronts it.
startsecs=5
stdout_logfile=/var/log/neurolink/cloudflared.log
stderr_logfile=/var/log/neurolink/cloudflared.err.log
priority=60
