# saas/vm-image/supervisor/programs/llama.conf
[program:llama]
command=/usr/local/bin/llama-server
    --model /models/Octen-Embedding-8B.f16.gguf
    --host 127.0.0.1
    --port 8081
    --embedding
    --ctx-size 4096
    --threads 4
autostart=true
autorestart=true
; Large model load can take 30-120s; give it room before supervisor gives up.
startsecs=30
startretries=5
stdout_logfile=/var/log/neurolink/llama.log
stderr_logfile=/var/log/neurolink/llama.err.log
priority=30
