# saas/vm-image/supervisor/programs/neo4j.conf
[program:neo4j]
command=/opt/neo4j/bin/neo4j console
directory=/opt/neo4j
autostart=true
autorestart=true
startsecs=15
stdout_logfile=/var/log/neurolink/neo4j.log
stderr_logfile=/var/log/neurolink/neo4j.err.log
environment=JAVA_OPTS="-Xms256m -Xmx1g"
priority=20
