# saas/vm-image/supervisor/programs/neuro-link.conf
[program:neuro-link]
command=/usr/local/bin/neuro-link serve --port 8080 --token-file /secrets/api-token
directory=/data
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/neurolink/neuro-link.log
stderr_logfile=/var/log/neurolink/neuro-link.err.log
environment=NLR_ROOT="/data",QDRANT_URL="http://127.0.0.1:6333",NEO4J_HTTP_URL="http://127.0.0.1:7474",LLAMA_URL="http://127.0.0.1:8081"
priority=40
