# saas/vm-image/supervisor/programs/qdrant.conf
[program:qdrant]
command=/usr/local/bin/qdrant --config-path /etc/qdrant/config.yaml
directory=/state/qdrant
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/neurolink/qdrant.log
stderr_logfile=/var/log/neurolink/qdrant.err.log
priority=10
